# systemd/scalping-bot.service
[Unit]
Description=Scalping Bot for Binance Futures
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/MRROBOT-FUTURE
Environment="PATH=/home/ubuntu/MRROBOT-FUTURE/venv/bin"
ExecStart=/home/ubuntu/MRROBOT-FUTURE/venv/bin/python -m src.main

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Logging
StandardOutput=append:/var/log/scalping-bot/output.log
StandardError=append:/var/log/scalping-bot/error.log
SyslogIdentifier=scalping-bot

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

# ============================================
# INSTRUÇÕES DE INSTALAÇÃO
# ============================================
#
# 1. Copie este arquivo para /etc/systemd/system/:
#    sudo cp systemd/scalping-bot.service /etc/systemd/system/
#
# 2. Ajuste o caminho e o usuário conforme sua configuração
#
# 3. Crie o diretório de logs:
#    sudo mkdir -p /var/log/scalping-bot
#    sudo chown ubuntu:ubuntu /var/log/scalping-bot
#
# 4. Recarregue o systemd:
#    sudo systemctl daemon-reload
#
# 5. Habilite o serviço para iniciar no boot:
#    sudo systemctl enable scalping-bot
#
# 6. Inicie o serviço:
#    sudo systemctl start scalping-bot
#
# 7. Verifique o status:
#    sudo systemctl status scalping-bot
#
# COMANDOS ÚTEIS:
# - Ver logs em tempo real: sudo journalctl -u scalping-bot -f
# - Reiniciar serviço: sudo systemctl restart scalping-bot
# - Parar serviço: sudo systemctl stop scalping-bot
# - Ver logs do arquivo: tail -f /var/log/scalping-bot/output.log
# ============================================
